# nginx.service
[Unit]
Description=nginx - high performance web server
Documentation=https://nginx.org/en/docs/
After=network-online.target remote-fs.target nss-lookup.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/opt/nginx/nginx.pid
ExecStartPre=/opt/nginx/bin/nginx -t -c /opt/nginx/nginx.conf
ExecStart=/opt/nginx/bin/nginx -c /opt/nginx/nginx.conf
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID



[Install]
WantedBy=multi-user.target
